# sshd/sshd.conf
[program:sshd]
autostart=true
command=/usr/sbin/sshd -D -e
stopsignal=INT
stdout_logfile=/tmp/sshd.log
